# deploy/local-zero/launchd/com.ashare.rt.eod.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ashare.rt.eod</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__ROOT__/deploy/local-zero/run_eod.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__ROOT__</string>
  <!-- Fire every 20 min; script only runs 15:20–18:30 Shanghai weekdays -->
  <key>StartInterval</key>
  <integer>1200</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>__ROOT__/deploy/local-zero/logs/eod.launchd.out</string>
  <key>StandardErrorPath</key>
  <string>__ROOT__/deploy/local-zero/logs/eod.launchd.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
